# systemd/clay-key-quota-retest.service
[Unit]
Description=Clay key quota retest (rapid retest of quota_exceeded rows)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=clay
WorkingDirectory=/opt/clay-enrichment
EnvironmentFile=/opt/clay-enrichment/.env
ExecStart=/opt/clay-enrichment/.venv/bin/python -m src.api_keys.cron.quota_retest
StandardOutput=journal
StandardError=journal
