# systemd/clay-key-revalidate.timer
[Unit]
Description=Daily revalidation of validated_keys

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
