# systemd/clay-key-revalidate.service
[Unit]
Description=Clay key revalidator (daily re-test of validated_keys)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=clay
WorkingDirectory=/opt/clay-enrichment
EnvironmentFile=/opt/clay-enrichment/.env
ExecStart=/opt/clay-enrichment/.venv/bin/python -m src.api_keys.cron.revalidate
StandardOutput=journal
StandardError=journal
