# deploy/clay-brief.service
[Unit]
Description=Clay enrichment campaign-brief HTTP service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=clay
WorkingDirectory=/opt/clay-enrichment
EnvironmentFile=/opt/clay-enrichment/.env
ExecStart=/opt/clay-enrichment/.venv/bin/python -m uvicorn src.web_api.app:app --host 127.0.0.1 --port 5001
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
